# sipXpublisher/bin/sipxpublisher.service.in
[Unit]
Description=sipxpublisher
Documentation=http://sipfoundry.org
After=network.target syslog.target

[Service]
Type=oneshot
ExecStart=@SIPX_SERVICEDIR@/sipxpublisher start
RemainAfterExit=true
ExecStop=@SIPX_SERVICEDIR@/sipxpublisher stop
StandardOutput=journal

[Install]
WantedBy=multi-user.target