# sipXregistry/bin/sipxregistrar.service.in
[Unit]
Description=sipxregistrar
Documentation=http://sipfoundry.org
After=network.target syslog.target

[Service]
Type=oneshot
ExecStart=@SIPX_SERVICEDIR@/sipxregistrar start
RemainAfterExit=true
ExecStop=@SIPX_SERVICEDIR@/sipxregistrar stop
StandardOutput=journal

[Install]
WantedBy=multi-user.target